# files/etc/uci-defaults/ovpn-dir.sh
#!/bin/sh

# Create folder for VPN files
mkdir /netaidkit/ovpn
chgrp www-data /netaidkit/ovpn
chmod g+w /netaidkit/ovpn
